# deploy/enactspace-mobile-money-reconcile.timer
[Unit]
Description=Run EnactSpace Mobile Money reconciliation every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=1min
Unit=enactspace-mobile-money-reconcile.service

[Install]
WantedBy=timers.target
